# packer/xcalar/rootfs/etc/systemd/system/xcalar.service.d/env.conf
[Service]
EnvironmentFile=-/var/lib/cloud/instance/ec2.env
